# factoriomcd.service
[Unit]
Description=Agent for FactorioMMO

[Service]
ExecStart=/opt/factoriommo-agent/env/bin/fmcd

[Install]
WantedBy=multi-user.target
